# supervisord.jenkins.conf
[program:jenkins]
command=/usr/bin/java
	-Dcom.sun.akuma.Daemon=daemonized
	-Djava.awt.headless=true
	-Djenkins.install.runSetupWizard=false
	-Djenkins.model.Jenkins.slaveAgentPort=%(ENV_JENKINS_SLAVE_PORT)s
	-jar
	/var/lib/jenkins/jenkins.war
	--debug=5
	--httpPort=-1
	--httpsKeyStore=%(ENV_JENKINS_HOME)s/jenkins.jks
	--httpsKeyStorePassword=%(ENV_JENKINS_KEYSTORE_PASSWORD)s
	--httpsPort=443
startsecs = 10
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
user = jenkins
